# linux/ollama-copilot.service.tmpl
[Unit]
Description=Ollama Copilot Service
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=/usr/local/bin/ollama-copilot -model=codellama:7b-code
Restart=on-failure
User={{USER}}

[Install]
WantedBy=multi-user.target
